# cf-under-attack.timer
[Unit]
Description=Run Cloudflare Under Attack Monitor every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s
RandomizedDelaySec=5s

[Install]
WantedBy=timers.target
